# obm/aua/bin/autoupdate-bsd
#!/bin/sh
#

# PROVIDE: @sed.script.name@
# REQUIRE: NETWORKING SERVERS DAEMON
# BEFORE: LOGIN
# KEYWORD: AutoUpdateAgent OBM

. /etc/rc.subr

name="@sed.script.name@"
rcvar=`set_rcvar`
extra_commands="usage"
start_cmd="@sed.script.name@_start"
stop_cmd="@sed.script.name@_stop"
usage_cmd="@sed.script.name@_usage"

@sed.script.name@_start()
{
    echo "Start AutoUpdateAgent for [ Ahsay Online Backup Manager ]$@"
    nohup sh -c "@sed.autoupdate.home@/bin/startup.sh" &> /dev/null
    sleep 5
}

@sed.script.name@_stop()
{
    echo "Shut down AutoUpdateAgent for [ Ahsay Online Backup Manager ]$@"
    nohup sh -c "@sed.autoupdate.home@/bin/shutdown.sh" &> /dev/null
    sleep 5
}

@sed.script.name@_usage()
{
    echo $"usage: $0 {start|stop}"
}

load_rc_config $name
run_rc_command "$1"